# mac/screentime/com.roblox.blocker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.roblox.blocker</string>

    <key>ProgramArguments</key>
    <array>
        <!-- path to python3 -->
        <string>/usr/bin/python3</string>
        <string>/Users/YOUR_USERNAME/path/to/roblox_blocker.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/roblox_blocker.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/roblox_blocker.err</string>
</dict>
</plist>